# com.davelydigest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.davelydigest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/daviduslan/Documents/DavelyDigest/run_local.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/daviduslan/Documents/DavelyDigest</string>

    <!-- Mon–Fri at 07:00 local time -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/daviduslan/Library/Logs/davelydigest.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/daviduslan/Library/Logs/davelydigest.error.log</string>

    <!-- Do not run immediately on load -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
